# debian/hello.postinst
#!/bin/bash

if [ -x "/etc/init.d/hello" ]; then
    update-rc.d hello defaults > /dev/null
    invoke-rc.d hello start || true
fi
